# files/wireguard_awgquick.in
#!/bin/sh

# PROVIDE: wireguard_amnezia
# REQUIRE: NETWORKING
# KEYWORD: shutdown
#
# wireguard_amnezia_enable (bool):    Set to "YES" to enable wireguard_amnezia.
#                             (default: "NO")
#
# wireguard_amnezia_interfaces (str): List of interfaces to bring up/down
#                             on start/stop. (eg: "wg0 wg1")
#                             (default: "")
# wireguard_amnezia_env (str):        Environment variables for the userspace
#                             implementation. (eg: "LOG_LEVEL=debug")

. /etc/rc.subr

name=wireguard_amnezia
rcvar=wireguard_amnezia_enable
extra_commands="reload status"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
reload_cmd="${name}_reload"
status_cmd="${name}_status"

wireguard_amnezia_start()
{
	${wireguard_amnezia_env:+eval export $wireguard_amnezia_env}

	for interface in ${wireguard_amnezia_interfaces}; do
		%%PREFIX%%/bin/awg-quick up ${interface}
	done
}

wireguard_amnezia_stop()
{
	for interface in ${wireguard_amnezia_interfaces}; do
		%%PREFIX%%/bin/awg-quick down ${interface}
	done
}

wireguard_amnezia_reload()
{
	${wireguard_amnezia_env:+eval export $wireguard_amnezia_env}

	for interface in ${wireguard_amnezia_interfaces}; do
		tmpfile="`mktemp`"
		%%PREFIX%%/bin/awg-quick strip ${interface} > ${tmpfile}
		%%PREFIX%%/bin/awg syncconf ${interface} ${tmpfile}
		rm -f ${tmpfile}
	done
}

wireguard_amnezia_status()
{
	${wireguard_amnezia_env:+eval export $wireguard_amnezia_env}
	wireguard_amnezia_status="0"

	for interface in ${wireguard_amnezia_interfaces}; do
		%%PREFIX%%/bin/awg show ${interface} || wireguard_amnezia_status="1"
	done

	return ${wireguard_amnezia_status}
}

load_rc_config $name

: ${wireguard_amnezia_enable="NO"}
: ${wireguard_amnezia_interfaces=""}
: ${wireguard_amnezia_env=""}

run_rc_command "$1"
